# stacks/sites-ginfo-apache/save-strategy/crontab
* * * * * root /etc/backup.sh
# Ne pas supprimer cette ligne, tout fichier cron doit se terminer par une ligne vide