# stacks/nextcloud/cron/crontab
*/5 * * * * php -f /var/www/html/nextcloud/cron.php
